# systemd/srt-init.service
[Unit]
Description=SRT HAL initialiser (runs srt-setup on start and on every RP2040 reconnect)
After=rp2040fs@srt.service
BindsTo=rp2040fs@srt.service

[Service]
Type=simple
User=srt
RuntimeDirectory=srt
RuntimeDirectoryMode=0755
ExecStart=/mnt/srt/hal/srt-init
Restart=always
RestartSec=5s
StartLimitIntervalSec=120
StartLimitBurst=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

