# systemd/srt-go.service
[Unit]
Description=SRT go/ axis controller (inotify mutual exclusion daemon)
After=rp2040fs@srt.service srt-init.service
BindsTo=rp2040fs@srt.service

[Service]
Type=notify
User=srt
ExecStartPre=/bin/sleep 2
ExecStart=/mnt/srt/srt-go
Restart=always
RestartSec=5s
StartLimitIntervalSec=120
StartLimitBurst=10
WatchdogSec=30s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
