# systemd/srt-go.service
[Unit]
Description=SRT HAL go/ axis controller (inotify mutual exclusion daemon)
After=rp2040fs@srt.service srt-init.service
BindsTo=rp2040fs@srt.service
Wants=srt-init.service

[Service]
Type=notify
User=srt
ExecStart=/mnt/srt/hal/srt-go
Restart=always
RestartSec=5s
StartLimitIntervalSec=120
StartLimitBurst=10
WatchdogSec=30s
NotifyAccess=all
TimeoutStartSec=60s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

