# systemd/srt-watchdog.service
[Unit]
Description=SRT watchdog — periodic health check and recovery
After=srt-go.service srt-gs232.service srt-init.service

[Service]
Type=oneshot
User=root
ExecStart=/mnt/srt/control/srt-watchdog
StandardOutput=journal
StandardError=journal

