# systemd/srt-watchdog.service
[Unit]
Description=SRT watchdog — periodic health check and recovery
After=srt-go.service srt-init.service

[Service]
Type=oneshot
# Runs as root — needs to call systemctl restart
User=root
ExecStart=/mnt/srt/srt-watchdog
StandardOutput=journal
StandardError=journal
