# systemd/rp2040fs@srt.service.d/allow_other.conf
[Service]
# Add -o allow_other so non-root users (e.g. srt) can read/write the FUSE mount.
# The blank ExecStart= clears the upstream value before setting the new one.
ExecStart=
ExecStart=/usr/local/bin/rp2040fs /mnt/rp2040 --device /dev/rp2040_gpio_fs -o allow_other

